# snapdSB.service
[Unit]
Description=Prepare mount points
Before=remote-fs-pre.target
Wants=remote-fs-pre.target

[Service]
Type=oneshot
ExecStart=bash /opt/snapdSB/snapdSB.sh
ExecStopPost=chattr +i /

[Install]
WantedBy=remote-fs.target
