# snapdSB.service
[Unit]
Description=Prepare mount points
Before=remote-fs-pre.target
Wants=remote-fs-pre.target

[Service]
Type=oneshot
ExecStartPre=chattr -i /
ExecStart=/bin/sh -c "([ -L '/snap' ] && echo 'snap symlink already exists') || \
(sudo ln -s '/var/mnt/snap' '/snap' && echo 'symlink /snap created')"
ExecStopPost=chattr +i /

[Install]
WantedBy=remote-fs.target
